# autostart/My-Peeps.desktop
[Desktop Entry]
Type=Application
Name=cibo
Exec=sh -c '/home/jules/git/My-Peeps/autostart/My-Peep.sh

